# gunicorn_TurtleMaster.socket
[Unit]
Description=TurtleMaster socket
[Socket]
ListenStream= /home/wwwadmin/TurtleMaster/TurtleMaster.sock 
[Install]
WantedBy=sockets.target
